# config/com.erd.embed-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.erd.embed-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/francesco/EF/ethereum-research-data/.venv/bin/erd-embed-proxy</string>
        <string>--port</string>
        <string>11435</string>
        <string>--max-chars</string>
        <string>4000</string>
        <string>--ollama-concurrency</string>
        <string>20</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/francesco/.cache/erd-embed-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/francesco/.cache/erd-embed-proxy.log</string>
</dict>
</plist>
