# config/com.erd.embed-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd service for the embed proxy (sits between Meilisearch and Ollama).
  Only needed for native (non-Docker) installs.

  Install:
    1. Edit paths below to match your checkout location
    2. cp config/com.erd.embed-proxy.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.erd.embed-proxy.plist

  The embed proxy is optional — it handles token-aware splitting for
  Meilisearch's built-in embedder. If you use batch_embed.py for
  embedding (recommended), the proxy is only needed for auto-embedding
  new documents added after the initial batch.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.erd.embed-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Edit this path to match your checkout -->
        <string>/path/to/ethereum-research-data/.venv/bin/erd-embed-proxy</string>
        <string>--port</string>
        <string>11435</string>
        <string>--max-chars</string>
        <string>4000</string>
        <string>--ollama-concurrency</string>
        <string>20</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/erd-embed-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/erd-embed-proxy.log</string>
</dict>
</plist>
